# suarm/tmpl/django_supervisor.conf
[program:{{project_name}}]
command = {{project_path}}/bin/start.sh
stdout_logfile = {{project_path}}/var/log/supervisor.log
redirect_stderr = true
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8